# services/nginx.service.conf
[program:nginx]
user=hpess
command=/usr/sbin/nginx -c /storage/nginx/nginx.conf
stdout_events_enabled = true
stderr_events_enabled = true
